# stuff/dipy_autosync.service
[Unit]
Description=My Script

[Service]
Type=forking
ExecStart=sh /home/pi/DiPy/update_loop.sh

[Install]
WantedBy=multi-user.target
